# nagios/nagios.ini
[program:nagios]
user=nagios
command=/usr/local/nagios/bin/nagios /usr/local/nagios/etc/nagios.cfg
stdout_logfile=/dev/stdout
redirect_stderr=true

[program:apahce2]
command=apachectl -e info -DFOREGROUND
stdout_logfile=/dev/stdout
redirect_stderr=true
